# supervisor.conf
[program:h2]
command=nohup /opt/h2/bin/h2-server.sh &
autorestart=true